# deploy/systemd/salsa.service.d/override.conf
[Service]
Environment="RUST_LOG=info"
Environment="KEY_FILE_PATH=<path>/privkey.pem"
Environment="CERT_FILE_PATH=<path>/fullchain.pem"
